# .config/systemd/user/gnome-keyring.service
[Unit]
Description=GNOME Keyring (secrets) for Hyprland session
Before=nextcloud.service

[Service]
Type=simple
ExecStart=/usr/bin/gnome-keyring-daemon --components=secrets --foreground
Restart=on-failure
Environment=XDG_RUNTIME_DIR=%t

[Install]
WantedBy=hyprland-session.target
