# .config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland session target. This is the root in the auto-start dependency tree. See hyprland.conf.
