# .config/systemd/user/nextcloud.service
[Unit]
Description=Nextcloud Client

[Service]
ExecStart=/usr/bin/nextcloud
Restart=on-failure

[Install]
WantedBy=hyprland-session.target
